# netdata.conf
[program:netdata]
command=/usr/sbin/netdata -nd -p ${NETDATA_PORT}
